# conf/etc/network/if-up.d/telecaster-rsync
#!/bin/sh

su -c "/usr/local/lib/telecaster/instance/manage.py telecaster-rsync-archives" telecaster

