# tcserver/conf/etc/cron.d/telecaster_cleanup
0 12 * * * telecaster python /home/telecaster/apps/tools/cleanup/rm_date_days.py 8 /home/telecaster/trash/ >> /home/telecaster/logs/cleanup.log
