# tcserver/conf/etc/cron.d/telecaster_rsync
0 */4 * * * telecaster /usr/local/lib/telecaster/instance/manage.py telecaster-rsync-archives
