# tcserver/conf/etc/cron.d/telecaster_restart
0 8 * * * root /usr/local/bin/telecaster-restart

