# tcserver/conf/etc/cron.d/telecaster_import
*/10 * * * * telecaster /usr/local/lib/telecaster/instance/manage.py teleforma-import-conferences Pre-Barreau ~/logs/import.log
