# tcserver/conf/etc/network/if-up.d/telecaster-rsync
#!/bin/sh

/usr/local/lib/telecaster/instance/manage.py teleforma-import-conferences Pre-Barreau /var/log/telecaster/import.log

