# conf/etc/cron.d/telecaster_rsync
15 12 * * * telecaster /usr/local/lib/telecaster/instance/manage.py telecaster-rsync-archives
15 15 * * * telecaster /usr/local/lib/telecaster/instance/manage.py telecaster-rsync-archives
0  21 * * * telecaster /usr/local/lib/telecaster/instance/manage.py telecaster-rsync-archives
